# Packaging/debian/pharos-mesh.service
[Unit]
Description=Pharos Mesh local-first replica
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pharos-mesh
Group=pharos-mesh
StateDirectory=pharos-mesh
StateDirectoryMode=0750
ExecStart=/usr/bin/pharos-mesh distributed sync-serve --data-dir /var/lib/pharos-mesh --relay production
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/pharos-mesh

[Install]
WantedBy=multi-user.target
